# deploy/systemd/tg-cli-bot.service
[Unit]
Description=Telegram CLI Gateway Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgcli
Group=tgcli
WorkingDirectory=/opt/tg-cli-gateway
EnvironmentFile=/opt/tg-cli-gateway/.env
ExecStart=/opt/tg-cli-gateway/.venv/bin/python -m app.main
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
